# deploy/proca-blue.service
# Systemd service for Proca (blue instance, port 4000)
#
# Install:
#   sudo cp deploy/proca-blue.service /etc/systemd/system/
#   sudo systemctl enable proca-blue && sudo systemctl start proca-blue
#
# For zero-downtime deploys, also install proca-green.service.
# The deploy-bluegreen script handles start/stop of each.

[Unit]
Description=Proca (blue instance, port 4000)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=app
Group=app
WorkingDirectory=/srv/proca/current

Environment=PORT=4000
Environment=RELEASE_TMP=/srv/proca/current/tmp

# Load secret env vars (DATABASE_URL, SECRET_KEY_BASE, etc.)
EnvironmentFile=/srv/proca/environment

ExecStart=/srv/proca/current/bin/proca start
ExecStop=/srv/proca/current/bin/proca stop

# Health-check after start
ExecStartPost=/bin/sh -c 'for i in $$(seq 1 30); do curl -sf http://127.0.0.1:4000/health >/dev/null && exit 0; sleep 1; done; exit 1'

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
